# dot_config/systemd/user/ibus-custom.service
[Unit]
PartOf=graphical-session.target

Description=IBus Service
After=graphical-session.target

[Service]
Type=exec
ExitType=cgroup
ExecStartPre=/usr/bin/sleep 7
ExecStart=/usr/bin/ibus-daemon -drx
Restart=no
TimeoutSec=15s
Slice=app.slice

[Install]
WantedBy=graphical-session.target
